# packages/migrestore/migrestore.service
[Unit]
Description=Balena Migration Restore Service

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/bin/migrestore.sh

[Install]
WantedBy=multi-user.target
