# packages/migsupervisor/migsupervisor.service
[Unit]
Description=Balena Migration Supervisor Service
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/migSupervisor.sh

[Install]
WantedBy=multi-user.target
