# packages/mig3gconn/mig3gconn.service
[Unit]
Description=3G connection service
After=network.target
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/carrierConnect.sh

[Install]
WantedBy=multi-user.target
